# Project/ansible/install_db/tasks/files/backup.service
[Unit]
Description=MySQL database backup service
[Service]
Type=oneshot
ExecStart=/root/mysql-backup.sh